# Library/LaunchAgents/com.user.frpc-tx.plist.tmpl
{{- $frp := include "frp/encrypted_private_machines.toml.asc" | decrypt | fromToml -}}
{{- $machineID := index . "frp_machine" | default "" -}}
{{- if not $machineID -}}
{{- fail "No frp_machine configured. Set frp_machine in ~/.config/chezmoi/chezmoi.toml [data] to one of frp/machines.toml machines." -}}
{{- end -}}
{{- $machine := index $frp.machines $machineID -}}
{{- if not $machine -}}
{{- fail (printf "No FRP machine config for frp_machine %q. Add [machines.%s] to the encrypted FRP machine inventory." $machineID $machineID) -}}
{{- end -}}
{{- $defaults := index $frp "defaults" -}}
{{- $homeDir := index $machine "home_dir" | default (index $defaults "home_dir") -}}
{{- $frpcBin := index $machine "frpc_bin" | default (index $defaults "frpc_bin") -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.frpc-tx</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ $frpcBin }}</string>
        <string>-c</string>
        <string>{{ $homeDir }}/frp/frpc-tx.toml</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{ $homeDir }}/frp</string>
    <key>StandardOutPath</key>
    <string>{{ $homeDir }}/frp/frpc-tx.log</string>
    <key>StandardErrorPath</key>
    <string>{{ $homeDir }}/frp/frpc-tx.err.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
